# mbsync/.config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=3m
OnUnitActiveSec=3m
Unit=mbsync.service

[Install]
WantedBy=timers.target
